# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:ckb-rpc-server]
command=/bin/sh -c 'exec /usr/local/bin/ckb-rpc-server --ckb-rpc ${CKB_RPC_URL:-http://127.0.0.1:8114}'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=RUST_LOG="info"

[program:ckb-docs-server]
command=/usr/local/bin/ckb-docs-server --docs-path /app/docs
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=RUST_LOG="info"

[program:ckb-tools-server]
command=/bin/sh -c 'exec /usr/local/bin/ckb-tools-server --ckb-rpc ${CKB_RPC_URL:-http://127.0.0.1:8114}'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=RUST_LOG="info"